# deployment/example-api.service
[Unit]
Description=Example API Service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user
ExecStart=/usr/local/bin/example-api
Restart=always
RestartSec=3
Environment="PORT=80"

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=example-api

[Install]
WantedBy=multi-user.target
